# frps/root/etc/s6-overlay/s6-rc.d/svc-frps/run
#!/usr/bin/with-contenv bash

s6-setuidgid abc \
    frps -c ${FRP_CONFIG_PATH}/config.toml
